# supervisor_redis.conf
[program:redis]
command=/usr/local/bin/redis-server --appendonly true
autorestart=true
startretries=3


